# services/com.coyodude.dplussim/run
#!/bin/sh
set -eu
if set -o pipefail >/dev/null 2>&1; then
    :
fi

resolve_path() {
    target="$1"
    if command -v realpath >/dev/null 2>&1; then
        realpath "$target"
    elif command -v readlink >/dev/null 2>&1; then
        readlink -f "$target"
    else
        python3 -c 'import os, sys; print(os.path.realpath(sys.argv[1]))' "$target"
    fi
}

if [ -n "${DPLUS_SIM_ROOT:-}" ]; then
    PROJECT_ROOT="$(resolve_path "${DPLUS_SIM_ROOT}")"
else
    SCRIPT_PATH="$(resolve_path "$0")"
    SERVICE_DIR="$(dirname "${SCRIPT_PATH}")"
    SERVICES_ROOT="$(dirname "${SERVICE_DIR}")"
    PROJECT_ROOT="$(dirname "${SERVICES_ROOT}")"
fi

export PYTHONUNBUFFERED=1
exec python3 "${PROJECT_ROOT}/src/dplus_sim.py" "$@"
