# services/com.coyodude.dplussim/run
#!/usr/bin/env bash
set -euo pipefail

resolve_path() {
    local target="$1"
    if command -v realpath >/dev/null 2>&1; then
        realpath -- "$target"
        return
    fi

    if command -v python3 >/dev/null 2>&1; then
        python3 -c 'import os, sys; print(os.path.realpath(sys.argv[1]))' "$target"
        return
    fi

    if command -v readlink >/dev/null 2>&1; then
        local resolved
        if resolved=$(readlink -f -- "$target" 2>/dev/null); then
            printf '%s\n' "$resolved"
            return
        fi
    fi

    echo "Unable to resolve path: $target" >&2
    return 1
}

if [[ -n "${DPLUS_SIM_ROOT:-}" ]]; then
    PROJECT_ROOT="$(resolve_path "${DPLUS_SIM_ROOT}")"
else
    SCRIPT_PATH="$(resolve_path "${BASH_SOURCE[0]}")"
    SERVICE_DIR="$(dirname "${SCRIPT_PATH}")"
    SERVICES_ROOT="$(dirname "${SERVICE_DIR}")"
    PROJECT_ROOT="$(dirname "${SERVICES_ROOT}")"
fi

export PYTHONUNBUFFERED=1
exec python3 "${PROJECT_ROOT}/src/dplus_sim.py" "$@"
